# deploy/runpod/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/liveai

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:avatar]
command=uvicorn services.avatar.app:app --host 0.0.0.0 --port 8001
directory=/workspace/live-ai
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/liveai/avatar.log
stderr_logfile=/var/log/liveai/avatar_err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1",MODEL_CACHE_DIR="/workspace/models",LIVEPORTRAIT_MODEL_PATH="/workspace/models/liveportrait",NVIDIA_VISIBLE_DEVICES="all"

[program:voice]
command=uvicorn services.voice.app:app --host 0.0.0.0 --port 8002
directory=/workspace/live-ai
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/liveai/voice.log
stderr_logfile=/var/log/liveai/voice_err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1",MODEL_CACHE_DIR="/workspace/models",FISH_SPEECH_MODEL_PATH="/workspace/models/fish-speech",NVIDIA_VISIBLE_DEVICES="all"

[program:expression]
command=uvicorn services.expression.app:app --host 0.0.0.0 --port 8003
directory=/workspace/live-ai
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/liveai/expression.log
stderr_logfile=/var/log/liveai/expression_err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1",MODEL_CACHE_DIR="/workspace/models",RMBG_MODEL="/workspace/models/rmbg",NVIDIA_VISIBLE_DEVICES="all"

[program:streaming]
command=uvicorn services.streaming.app:app --host 0.0.0.0 --port 8004
directory=/workspace/live-ai
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/liveai/streaming.log
stderr_logfile=/var/log/liveai/streaming_err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1",MODEL_CACHE_DIR="/workspace/models",USE_NVENC="true",NVIDIA_VISIBLE_DEVICES="all"
